# bootstrap/files/archaudit/archaudit.timer
[Unit]
Description=Run archaudit.sh daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
