# bootstrap/files/fortune-motd.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "{{ fortune_bin }} -c -a archer buffy discworld firefly futurama hitchhiker > /etc/motd"
