# bootstrap/files/fortune-motd.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "fortune -c -a archer buffy discworld firefly futurama hitchhiker > /etc/motd"
